# dot_config/systemd/user/pam.service.tmpl
{{- if eq .chezmoi.os "linux" -}}
[Unit]
Description=pam (MCP gateway) — Personal Agent Multiplexer
Documentation=https://git.dreamhost.com/dreamhost/ai/pam
After=default.target

[Service]
Type=simple
Environment=PATH=%h/.local/share/mise/shims:%h/.local/bin:/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/usr/local/bin:/usr/bin:/bin
ExecStart={{ lookPath "pam" }} start --foreground
ExecReload={{ lookPath "pam" }} reload
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
{{- end }}
